# jupyterlab-base/profile.d/local01-cd.sh
#!/bin/sh
if [ -n "${HOME}" ];
   cd "${HOME}"
fi
